# setup_files/yt_client.service
[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/Desktop
ExecStart=python hmi

[Install]
WantedBy=multi-user.target
