# setup_files/yt_client.service
[Service]
Type=simple
User=max
WorkingDirectory=/home/max
ExecStart=/home/max/hmi/venv/bin/python hmi

[Install]
WantedBy=multi-user.target
